# files/etc/systemd/system/nvfancontrol.service
[Unit]
Description=Nvidia fan control

[Service]
ExecStart=/usr/bin/nvfancontrol
Restart=always

[Install]
WantedBy=default.target
